# systemd/evohome-logger.service
[Unit]
Description=Evohome logger (Podman) run-once
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/evohome-logger
Environment="PATH=/usr/local/bin:/usr/bin"
Environment="REPO_DIR=/workspace"
ExecStart=/usr/bin/env bash -c 'cd "${REPO_DIR:-/workspace}" && exec make run-once'

[Install]
WantedBy=multi-user.target
